# .config/profile.d/50-nodejs.sh
#!/bin/sh

# Node.js
export NODE_REPL_HISTORY="$XDG_DATA_HOME/node/repl_history"

# NPM
export NPM_CONFIG_USERCONFIG="$XDG_CONFIG_HOME/npm/config"
export NPM_CONFIG_CACHE="$XDG_CACHE_HOME/npm"

